# sudo-logger-agent.service
[Unit]
Description=sudo-logger audit agent (plugin handler + eBPF recorder)
# Must start after systemd-logind (cgroup session scopes) and network (logserver).
After=systemd-logind.service network-online.target
Wants=network-online.target
# Prevent manual stop to close the attack vector where a user stops the agent
# to avoid logging.  Any kill -9 typed inside a sudo session is itself captured
# as TTY I/O before the agent dies, making the tampering self-documenting.
RefuseManualStop=yes

[Service]
ExecStart=/usr/bin/sudo-logger-agent -config /etc/sudo-logger/agent.conf
KillMode=process
Restart=always
RestartSec=2

User=root
# The agent runs as root with full capabilities (CAP_BPF, CAP_SYS_PTRACE,
# CAP_SYS_ADMIN, etc.) because eBPF tracing and the BPF LSM sandbox require
# them. Systemd sandboxing flags that create private mount namespaces or add
# seccomp/Landlock filters (PrivateTmp, PrivateDevices, ProtectKernelTunables,
# ProtectKernelModules, LockPersonality, RestrictSUIDSGID) cause inotify_add_watch
# to return EACCES for /etc and /proc/sys paths inside the service's namespace,
# preventing the sandbox watcher from tracking atomic file replacements.
# The BPF LSM sandbox (kernel-level) is the meaningful confinement boundary.
#
# BPF program pinning requires write access to /sys/fs/bpf.
ReadWritePaths=/sys/fs/bpf

[Install]
WantedBy=multi-user.target
